# upstart/event-injectors.conf
description "start up the injector instances on boot"

start on runlevel [2345]

task

script
  start event-injector QUEUE=events
  start event-injector QUEUE=errors
end script
